# env/flow/flow_dumper.ini
[program:flow_dumper]
command=/usr/bin/python2 flow_dumper.py
numprocs=1
process_name=%(program_name)s
directory=/root
user=root
autostart=true
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stdout_logfile=/var/log/flow_dumper_out.log
loglevel=info

[program:container_mon]
command=/usr/bin/python2 container_mon.py
numprocs=1
process_name=%(program_name)s
directory=/root
user=root
autostart=true
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stdout_logfile=/var/log/container_mon_out.log
loglevel=info